# snapshot-fix.service
[Unit]
Description=Claude Shell Snapshot Auto-Fixer
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/work/tmux-orchestrator/snapshot-monitor.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target